# udev/70-snap.nmap.rules
# This file is automatically generated.
# network-control
KERNEL=="rfkill", TAG+="snap_nmap_ncat"
# network-control
KERNEL=="tun", TAG+="snap_nmap_ncat"
TAG=="snap_nmap_ncat", SUBSYSTEM!="module", SUBSYSTEM!="subsystem", RUN+="/usr/lib/snapd/snap-device-helper $env{ACTION} snap_nmap_ncat $devpath $major:$minor"
# network-control
KERNEL=="rfkill", TAG+="snap_nmap_nmap"
# network-control
KERNEL=="tun", TAG+="snap_nmap_nmap"
TAG=="snap_nmap_nmap", SUBSYSTEM!="module", SUBSYSTEM!="subsystem", RUN+="/usr/lib/snapd/snap-device-helper $env{ACTION} snap_nmap_nmap $devpath $major:$minor"
# network-control
KERNEL=="rfkill", TAG+="snap_nmap_nping"
# network-control
KERNEL=="tun", TAG+="snap_nmap_nping"
TAG=="snap_nmap_nping", SUBSYSTEM!="module", SUBSYSTEM!="subsystem", RUN+="/usr/lib/snapd/snap-device-helper $env{ACTION} snap_nmap_nping $devpath $major:$minor"
