# udev/70-snap.chirp-snap.rules
# This file is automatically generated.
# opengl
KERNEL=="dma_buf_te", TAG+="snap_chirp-snap_chirp"
# opengl
KERNEL=="galcore", TAG+="snap_chirp-snap_chirp"
# opengl
KERNEL=="mali[0-9]*", TAG+="snap_chirp-snap_chirp"
# opengl
KERNEL=="nvhost-*", TAG+="snap_chirp-snap_chirp"
# opengl
KERNEL=="nvmap", TAG+="snap_chirp-snap_chirp"
# opengl
KERNEL=="pvr_sync", TAG+="snap_chirp-snap_chirp"
# opengl
KERNEL=="renderD[0-9]*", TAG+="snap_chirp-snap_chirp"
# opengl
KERNEL=="tegra_dc_[0-9]*", TAG+="snap_chirp-snap_chirp"
# opengl
KERNEL=="tegra_dc_ctrl", TAG+="snap_chirp-snap_chirp"
# opengl
KERNEL=="vchiq", TAG+="snap_chirp-snap_chirp"
# opengl
KERNEL=="vcsm-cma", TAG+="snap_chirp-snap_chirp"
# opengl
SUBSYSTEM=="dma_heap", KERNEL=="linux,cma", TAG+="snap_chirp-snap_chirp"
# opengl
SUBSYSTEM=="dma_heap", KERNEL=="system", TAG+="snap_chirp-snap_chirp"
# opengl
SUBSYSTEM=="drm", KERNEL=="card[0-9]*", TAG+="snap_chirp-snap_chirp"
# opengl
SUBSYSTEM=="host1x-fence", TAG+="snap_chirp-snap_chirp"
# opengl
SUBSYSTEM=="nvidia-gpu-v2" KERNEL=="ctrl", TAG+="snap_chirp-snap_chirp"
# opengl
SUBSYSTEM=="nvidia-gpu-v2" KERNEL=="prof", TAG+="snap_chirp-snap_chirp"
# opengl
SUBSYSTEM=="nvidia-gpu-v2-power" KERNEL=="power", TAG+="snap_chirp-snap_chirp"
# raw-usb
SUBSYSTEM=="tty", ENV{ID_BUS}=="usb", TAG+="snap_chirp-snap_chirp"
# raw-usb
SUBSYSTEM=="usb", TAG+="snap_chirp-snap_chirp"
# raw-usb
SUBSYSTEM=="usbmisc", TAG+="snap_chirp-snap_chirp"
TAG=="snap_chirp-snap_chirp", SUBSYSTEM!="module", SUBSYSTEM!="subsystem", RUN+="/usr/lib/snapd/snap-device-helper $env{ACTION} snap_chirp-snap_chirp $devpath $major:$minor"
